# src/profile.d/00-ps4.sh
# shellcheck shell=bash

_ps4
PS4='(${BASH_SOURCE}:${LINENO}): ${FUNCNAME[0]} - [${SHLVL},${BASH_SUBSHELL}, $?]'
