# src/profile.d/node.sh
# shellcheck shell=sh

# Node default package manager
export NODE_COREPACK='pnpm'

# Node prefix
export NODE_PREFIX='/node'
