# src/profile.d/password.sh
# shellcheck shell=sh

# Admin password
export PASSWORD=''
