# rc/profile.d/00-python.sh
# shellcheck shell=sh

# Python' program name
export PYTHON_NAME='py'

# Python' prefix for self packages
export PYTHON_PREFIX="/${PYTHON_NAME}"
alias py='cd "${PYTHON_PREFIX}"'

# Python' version (major.minor) for self packages
export PYTHON_VERSION='3.10'

