# backup/src/profile.d/00-go.sh
# shellcheck shell=sh

# GOHOME
export GOHOME='/go'
